# a3mix/raspberry/etc/systemd/system/a3_mic.service
[Unit]
Description=A3 Mixer Controller
After=network.target

[Service]
#User=aaa
Environment="PATH=/home/aaa/.local/bin:$PATH"
ExecStart=/usr/bin/python3 /home/aaa/a3-system/a3mix/python/mic.py

[Install]
WantedBy=multi-user.target
